# deploy/pi/gabiru-update.service
[Unit]
Description=Gabiru auto-update (git pull + restart)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/env bash /opt/gabiru/deploy/pi/gabiru-update.sh
