# deploy/pi/gabiru.service
[Unit]
Description=Gabiru 3D Printer Manager
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/gabiru/backend
ExecStart=/opt/gabiru/backend/.venv/bin/python -m uvicorn app.main:app --app-dir /opt/gabiru/backend --host 0.0.0.0 --port 8080
Restart=always
RestartSec=2

# Optional: autoconnect on boot
Environment=GABIRU_AUTOCONNECT=1
# Use the fixed USB ACM port + 115200 baud.
Environment=GABIRU_PORT=/dev/ttyACM0
Environment=GABIRU_BAUDRATE=115200

# Timelapse mode: 'interval' (default) or 'layer' (OctoPrint-like).
# Environment=GABIRU_TIMELAPSE_MODE=layer

[Install]
WantedBy=multi-user.target
